# systemd/sonarr.service
[Unit]
Description=Sonarr Container (via WireGuard)
After=wireguard.service
Requires=wireguard.service
BindsTo=wireguard.service
PartOf=wireguard.service

[Service]
Type=simple
TimeoutStartSec=0
Restart=on-failure
RestartSec=10

# Wait for WireGuard tunnel to be established
ExecStartPre=/home/tramfjord/projects/mediaserver/systemd/wait-for-wireguard.sh

# Remove existing container if it exists
ExecStartPre=-/usr/bin/docker stop sonarr
ExecStartPre=-/usr/bin/docker rm sonarr
ExecStartPre=/usr/bin/docker pull lscr.io/linuxserver/sonarr:latest

ExecStart=/usr/bin/docker run \
    --name sonarr \
    --network=container:wireguard \
    -e PUID=1000 \
    -e PGID=1000 \
    -e TZ=Etc/UTC \
    -v /home/tramfjord/projects/mediaserver/config/sonarr:/config \
    -v /data:/data \
    --log-driver=journald \
    lscr.io/linuxserver/sonarr:latest

ExecStop=/usr/bin/docker stop sonarr
ExecStopPost=-/usr/bin/docker rm sonarr

[Install]
WantedBy=multi-user.target
