# systemd/wireguard.service.erb
[Unit]
Description=WireGuard VPN Container
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
TimeoutStartSec=0
Restart=on-failure
RestartSec=10

ExecStart=/usr/bin/docker compose -f <%= compose_file %> up wireguard
ExecStop=/usr/bin/docker compose -f <%= compose_file %> stop wireguard

[Install]
WantedBy=multi-user.target
