# systemd/wireguard.service
[Unit]
Description=WireGuard VPN Container
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
TimeoutStartSec=0
Restart=on-failure
RestartSec=10

# Remove existing container if it exists
ExecStartPre=-/usr/bin/docker stop wireguard
ExecStartPre=-/usr/bin/docker rm wireguard
ExecStartPre=/usr/bin/docker pull lscr.io/linuxserver/wireguard:latest

ExecStart=/usr/bin/docker run \
    --name wireguard \
    --cap-add=NET_ADMIN \
    --cap-add=SYS_MODULE \
    --sysctl=net.ipv4.conf.all.src_valid_mark=1 \
    --sysctl=net.ipv6.conf.all.disable_ipv6=1 \
    -e PUID=1000 \
    -e PGID=1000 \
    -e TZ=Etc/UTC \
    -v /home/tramfjord/projects/mediaserver/config/wireguard:/config \
    -v /lib/modules:/lib/modules:ro \
    -p 7878:7878 \
    -p 8989:8989 \
    -p 9696:9696 \
    --log-driver=journald \
    lscr.io/linuxserver/wireguard:latest

ExecStop=/usr/bin/docker stop wireguard
ExecStopPost=-/usr/bin/docker rm wireguard

[Install]
WantedBy=multi-user.target
